# dist/lib/systemd/system/quad.service
[Unit]
Description=Quad Discord Bot
After=network.target

[Service]
Type=simple
ExecStart=quad
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
